# vps-agent/agunk-agent.service
[Unit]
Description=Agunk VPS Monitoring Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/agunk-agent.env
WorkingDirectory=/opt/agunk-agent
ExecStart=/opt/agunk-agent/.venv/bin/python -m uvicorn agunk_agent:app --host ${AGUNK_HOST} --port ${AGUNK_PORT}
Restart=always
RestartSec=3
User=root

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
