# vps-agent/sontoloyo-agent.service
[Unit]
Description=PT Sontoloyo Monitor — VPS Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Force the agent process to compute "today" in Asia/Jakarta (WIB),
# regardless of the host's `timedatectl` setting. This matters for
# `_count_active_ssh_in_db` — it compares a parsed expiry date against
# `datetime.now()` to decide whether a subscriber is still active.
# A VPS provisioned in Singapore, Tokyo, or Frankfurt would otherwise
# see "today" shift by 1-7 hours from the operator's mental model and
# either over- or under-count expiring accounts at the day boundary.
# Operators that need a different TZ can override this in the env file.
Environment="TZ=Asia/Jakarta"
EnvironmentFile=/etc/sontoloyo-agent.env
WorkingDirectory=/opt/sontoloyo-agent
ExecStart=/opt/sontoloyo-agent/.venv/bin/python -m uvicorn sontoloyo_agent:app --host ${SONTOLOYO_HOST} --port ${SONTOLOYO_PORT}
Restart=always
RestartSec=3
User=root

# Persistent state directory — used by the daily Ookla speedtest
# scheduler to cache the last result across agent restarts. systemd
# ensures /var/lib/sontoloyo exists with mode 0755 owned by root.
StateDirectory=sontoloyo
StateDirectoryMode=0755

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
